# rsync/templates/supervisor_filebeat.ini
[program:filebeat]
command=/start_filebeat.sh
process_name=%(program_name)s
directory=/home/filebeat
autostart = true
autorestart = unexpected
priority = 10
stdout_logfile=/var/log/supervisor/filebeat.log
redirect_stderr=true
stopsignal=INT
exitcodes=0,3
startsecs=0
