# rsync/templates/supervisor_cron.ini
[program:cron]
command = crond -f
startsecs = 5
stopwaitsecs = 3600
stopasgroup = false
killasgroup = true
